# roles/exporters/templates/gpu_process_exporter.service.j2
[Unit]
Description=GPU Process Exporter
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ ansible_user }}
Group={{ ansible_user }}
ExecStart={{ base_dir }}/exporters/bin/gpu_process_exporter
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gpu-process-exporter
Environment=RUST_LOG=info
Environment=BIND_ADDRESS=0.0.0.0:{{ exporter_ports.gpu }}

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ base_dir }}/exporters /var/log/cmstack

# Allow access to GPU devices and system info
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidiactl rw
DeviceAllow=/dev/nvidia-uvm rw
SupplementaryGroups=video

[Install]
WantedBy=multi-user.target
